# networking/crontab/ping1.sh
#!/bin/bash
/usr/bin/python3 /home/cwc/bash/ping-subprocess.py  >> /home/cwc/https/logs/ping1.txt
